# mm_moudles/mm_grpc_app/mm_grpc_app.ini
[program:mm_grpc_app]
directory=/usr/local/services/mm_grpc_app/
command=node mm_grpc_app.js
autostart=true
autorestart=true
startsecs=10
startretries=100